# filex.service
[Unit]
Description=FileX — file server with Markdown rendering
After=network.target

[Service]
Type=simple
ExecStart=%h/code/filex/dist/filex --root %h --port 9090 --bind 0.0.0.0
Restart=always
RestartSec=5
WorkingDirectory=%h

[Install]
WantedBy=default.target
