# filex.service
[Unit]
Description=FileX — file server with Markdown rendering
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/filex --root %h --port 9090 --bind 0.0.0.0
Restart=always
RestartSec=5
WorkingDirectory=%h
Environment=FILEX_BIN=%h/.local/bin/filex

[Install]
WantedBy=default.target
